# rootfs/data/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
redirect_stderr=true

[program:laminar]
command=laminard
directory=/data
killasgroup=true
stopasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:webhook]
command=webhook -port 9001 -hotreload -hooks /data/cfg/webhook.json -verbose
directory=/data
killasgroup=true
stopasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
